# docker-image-cleanup.timer
[Unit]
Description=Clean up docker images daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target

