# docker-container-cleanup.service
[Unit]
Description=Clean up old docker containers

[Service]
Type=oneshot
ExecStart=/usr/bin/docker-container-cleanup

