# docker-container-cleanup.timer
[Unit]
Description=Clean up docker containers daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target

